# Backend/deploy/hp-life-tracker-bot.service
[Unit]
Description=HP Life Tracker Telegram Bot
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/hp-life-tracker
EnvironmentFile=/opt/hp-life-tracker/.env
ExecStart=/usr/bin/docker compose exec -T app python run_telegram_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
